# modules/route53/templates/r53.erb
#!/bin/sh
#
# dns_update updates the cname record for this host in Amazon route 53
#
#       Written by Ken Garland <garlandkr@gmail.com>
#

### BEGIN INIT INFO
# Provides:          dns_update
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: dns_update for AWS route 53
# Description:       A script to update the dns entry of this box on boot
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
NAME=dns_update
DESC="dns update"

export AWS_ACCESS_KEY_ID=<%= scope.lookupvar('route53::params::access_key') %>
export AWS_SECRET_ACCESS_KEY=<%= scope.lookupvar('route53::params::secret_access_key') %>

prog=$(basename $0)
logger="logger -t $prog"

$logger "Running Route53"

cliscript=<%= scope.lookupvar('route53::params::cliscript') %>
zone=<%= scope.lookupvar('route53::params::domain') %>
hostname=<%= scope.lookupvar('route53::params::hostname') %>
type=<%= scope.lookupvar('route53::params::type') %>
address=<%= scope.lookupvar('route53::params::address') %>
ttl=<%= scope.lookupvar('route53::params::ttl') %>

route53=$($cliscript rrcreate $zone $hostname $type $address --ttl $ttl --replace)

$logger $route53
